# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

; FastAPI: webhook receiver + manual sync + health check.
; Only enqueues Celery tasks — never does the actual Outline/Qdrant work
; itself, so it can restart freely without losing in-flight work.
[program:fastapi]
command=uvicorn indexer.main:app --host 0.0.0.0 --port 8000
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; Celery worker: does the actual webhook/sync processing (chunk, embed, upsert).
[program:celery_worker]
command=celery -A indexer.celery_app worker --loglevel=info
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; Celery beat: enqueues a periodic incremental sync (SYNC_INTERVAL_SECONDS).
[program:celery_beat]
command=celery -A indexer.celery_app beat --loglevel=info
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; MCP server: search_knowledge tool over Streamable HTTP for AI agent clients.
[program:mcp_server]
command=python -m mcpserver.main
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
